# Sendspin/rootfs/etc/s6-overlay/s6-rc.d/init/run
#!/usr/bin/with-contenv bashio
    if [ ! -f "/config/settings-daemon.json" ]; then
        cp -p -r "/root/.config/sendspin.default/"* "/config"
        bashio::log.notice "NO SENDSPIN CONFIGURATION-FILES FOUND, RESTORED DEFAULTS!"
    fi
    if [ "$(bashio::config 'volume_controls')" = "Boosted" ]; then amixer -q set Master 100%; bashio::log.notice "VOLUME SET TO 100%!";fi
    ln -s "/config/" "/root/.config/sendspin"
    bashio::log.info "FINISHED INITIALIZATION!"