# Sendspin/rootfs/etc/s6-overlay/s6-rc.d/app/run
#!/command/with-contenv bashio




#apk info | sort

    function check_config {
        VALUE="$(bashio::config "$1")"
        if [[ "${1}" == str_* ]] && [ "${VALUE}" != null ]; then
            declare -g SENDSPIN_CMD+=" ${2} \"${VALUE}\""
        elif [[ "${1}" == bool_* ]]; then
            if [ "${VALUE}" != null ]; then
                declare -g "$1"="$VALUE"
            else
                declare -g "$1"="false"
            fi
        fi
    }


# CHECK CONFIGFILE
    bashio::log.info "Preparing Sendspin ..."
    #ln -s "/config/" "/root/.config/sendspin"
    ln -s "/config/" "/root/.config/sendspin"
    if [ ! -f "/config/settings-daemon.json" ]; then
        bashio::log.warning "Creating missing \"settings-daemon.json\" ..."
        { printf '%s\n' \
            "{"                                                   \
            "  \"name\": \"Home Assistant\","                     \
            "  \"client_id\": \"home_assistant\","                \
            "  \"last_server_url\": \"ws://[::1]:8927/sendspin\"" \
            "}"
        } > "/config/settings-daemon.json";
    fi
ls /root/.config/sendspin

# LOAD APP-CONFIG
    declare -g SENDSPIN_CMD="sendspin daemon --disable-mpris --hardware-volume false --manufacturer \"Music Assistant\" --product-name \"Home Assistant\""
    check_config bool_boost_volume   false
    check_config str_name            "--name"
    check_config str_id              "--id"
    check_config str_url             "--url"
    check_config str_log_level       "--log-level"
    #jq -r 'to_entries[] |  "\(.key): \(.value)"' "/config/settings-daemon.json"
    if [ "${bool_boost_volume}" = true ]; then bashio::log.info "Setting volume to 100% ..."; amixer -q set Master 100%; fi

    bashio::log.info "Starting Sendspin ..."
    eval "${SENDSPIN_CMD}"
    bashio::log.info "Stopping Sendspin ..."