# LedFx/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
sleep 2
bashio::log.info "Starting Nginx ..."
nginx -g "daemon off;"