# CUPS/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
    nginx -g "daemon off;" # > /dev/null 2>&1 &

# OVERLAY
    pid=$!
    while [ ! -S /run/nginx.sock ]; do sleep 0.1; done
    bashio::log.info "STARTED NGINX ..."
    printf '\n' >&3
    wait "$pid"