# CUPS/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
    nginx -g "daemon off;" &

# OVERLAY
    pid=$!
    while ! nc -z 127.0.0.1 630; do sleep 0.1; done
    bashio::log.info "STARTED NGINX!"
    printf '\n' >&3
    wait "$pid"