# CUPS/rootfs/etc/s6-overlay/s6-rc.d/dbus/run
#!/usr/bin/with-contenv bashio
    dbus-daemon --system --nofork > /dev/null 2>&1 &

# OVERLAY
    pid=$!
    while [ ! -S /run/dbus/system_bus_socket ]; do sleep 0.1; done
    bashio::log.info "STARTED DBUS ..."
    printf '\n' >&3
    wait "$pid"