# LedFx/rootfs/etc/s6-overlay/s6-rc.d/ledfx/run
#!/command/with-contenv bashio
    function check_config_str {
        declare -g "cfg_$1"="$(bashio::config "$1")"
    }


# APP-CONFIG
    bashio::log.info "Preparing LedFx ..."
    check_config_str log_level
    if [ "${cfg_log_level?}" = "INFO" ]; then
        declare -g cmd_log_level="-v"
    elif [ "${cfg_log_level?}" = "DEBUG" ]; then
        declare -g cmd_log_level="-vv"
    else
        declare -g cmd_log_level=""
    fi


# CONFIG-FILE
    if [ ! -d "/root/.ledfx" ]; then ln -s "/config/" "/root/.ledfx"; fi
    if [ ! -f "/config/config.json" ]; then { printf '%s\n' \
        "{" \
        "  \"audio\": {" \
        "    \"audio_device\": 2" \
        "  }," \
        "  \"configuration_version\": \"2.3.6\"," \
        "  \"sendspin_servers\": {" \
        "    \"music-assistant\": {" \
        "      \"client_name\": \"LedFx\"," \
        "      \"server_url\": \"ws://[::1]:8927/sendspin\"" \
        "    }" \
        "  }" \
        "}"
    } > "/config/config.json" ; fi


# START
    bashio::log.info "Starting LedFx ..."
    ledfx ${cmd_log_level} --offline --no-tray --host :: --port 8888