# LedFx/rootfs/etc/s6-overlay/s6-rc.d/ledfx/run
#!/command/with-contenv bashio
    declare -g "cfg_log_level"="$(bashio::config "log_level")"
    if   [ "${cfg_log_level?}" = "Debug" ]; then
        ledfx --offline --no-tray --host :: --port 8888 -vv &
    elif [ "${cfg_log_level?}" = "Info" ]; then
        ledfx --offline --no-tray --host :: --port 8888 -v &
    else
        ledfx --offline --no-tray --host :: --port 8888 &
    fi

# OVERLAY
    pid=$!
    while ! nc -z ::1 8888; do sleep 0.1; done
    bashio::log.info "STARTED LEDFX!"
    printf '\n' >&3
    wait "$pid"


    
